# tools/scripts/cosmic-clipboard.service
[Unit]
Description=Cosmic Clipboard History Daemon
After=graphical-session.target

[Service]
Type=simple
# We use a shell wrapper to start both text and image watchers
ExecStart=/bin/bash -c "wl-paste --type text --watch cliphist store & wl-paste --type image --watch cliphist store"
Restart=always
RestartSec=5

[Install]
WantedBy=cosmic-session.target
